# packages/device-cloud/device_cloud-18.6.19.tar.gz/device_cloud-18.6.19/share/device-cloud-watchdog.service.in
[Unit]
Description=WindRiver IoT Service Watchdog

[Service]
Type=forking
ExecStart=-%bindir%/device-cloud-watchdog start
ExecStop=-%bindir%/device-cloud-watchdog stop

[Install]
WantedBy=multi-user.target
