# packages/simplegitweb/simplegitweb-0.1.6.tar.gz/simplegitweb-0.1.6/simplegitweb/systemd/simplegitweb.service
[Unit]
Description=Simplegitweb Service
After=network.target
Wants=network.target

[Service]
User=simplegitweb
Group=simplegitweb
Type=simple
ExecStart=/usr/bin/gitwebserver.py --config /etc/simplegitweb.conf
Restart=on-failure

[Install]
WantedBy=multi-user.target