# packages/device-cloud/device_cloud-18.6.19-py3-none-any.whl/device_cloud-18.6.19.data/data/share/device_cloud/share/device-manager-sysv-init.in
#!/bin/sh
### BEGIN INIT INFO
# Provides:		sshd
# Required-Start:	$remote_fs $syslog $networking
# Required-Stop:	$remote_fs $syslog
# Default-Start:	2 3 4 5
# Default-Stop:		1
# Short-Description:	Device Manager App
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
NAME=device_manager.py
DAEMON_INTERPRETER=python
DAEMON_CMD="${DAEMON_INTERPRETER} /usr/bin/${NAME}"
DAEMON_OPTS="--log_to_syslog"
DESC="Device manager"
PIDFILE=/var/run/device-manager.pid
USER=%user%

set -e

do_start(){
	echo -n "Starting $DESC: "
	start-stop-daemon  --start --background --user $USER --chuid $USER -p $PIDFILE \
	  --make-pidfile --startas $DAEMON_CMD -- $DAEMON_OPTS
	echo "Status $?"
}

do_stop(){
	echo -n "Stopping $DESC: "
	start-stop-daemon -K -x "$DAEMON_INTERPRETER" -p $PIDFILE
	echo "Status $?"
}

case "$1" in
  start|stop)
	do_${1}
	;;
  restart|force-reload)
	do_stop
	sleep 1
	do_start
	;;
  *)
	N=/etc/init.d/$NAME
	echo "Usage: $N {start|stop|restart|force-reload}" >&2
	exit 1
	;;
esac

exit 0
